# deploy/yoonassist-backend.service
[Unit]
Description=YoonAssist AI Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=yoonassist
Group=yoonassist
WorkingDirectory=/opt/yoonassist
Environment="PATH=/opt/yoonassist/venv/bin"
EnvironmentFile=/opt/yoonassist/.env
ExecStart=/opt/yoonassist/venv/bin/uvicorn src.server:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10

# Limites de ressources
MemoryMax=800M
CPUQuota=80%

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yoonassist-backend

[Install]
WantedBy=multi-user.target

